# clients/wickrio/platforms/linux/scripts/WickrIOSvr.service
[Unit]
Description=WickrIO Client server

[Service]
ExecStart=/usr/sbin/WickrIOSvrINSTALLEXT
ExecStop=/usr/bin/wickrio_service_controlEXTENSION WickrIOSvrINSTALLEXT stop
