# clients/wickrio/platforms/linux/scripts/WickrIOCSvr.service
[Unit]
Description=WickrIO Console server

[Service]
ExecStart=/usr/sbin/WickrIOCSvrINSTALLEXT
