# clients/wickrio/platforms/linux/scripts/WickrIOCSvr.service
[Unit]
Description=WickrIO Console server

[Service]
ExecStart=/usr/sbin/WickrIOCSvrINSTALLEXT
ExecStop=/usr/bin/wickrio_service_controlEXTENSION WickrIOCSvrINSTALLEXT stop
